# containers/pootle/supervisor/nginx.conf
[program:nginx]
command=nginx -g "daemon off;"
autorestart=true
