# docker/pootle/supervisor/pootle.ini
[program:rqworker]
command=pootle rqworker
autorestart=true

[program:server]
command=pootle start
autorestart=true
